# scripts/raspArch/rk_autoStart.service
[Unit]
Description=AUTOSTART

[Service]
ExecStart=/bin/rk_autoStart

[Install]
WantedBy=multi-user.target
