# scripts/raspArch/rk_monitorTorrent.service
[Unit]
Description=Monitorar os arquivos transferidos por ftp

[Service]
ExecStart=/bin/rk_monitorTorrent

[Install]
WantedBy=multi-user.target
