# scripts/Linux/raspArch/old/rk_conexoesWifi.service
[Unit]
Description=monitorar as redes wifi disponiveis

[Service]
ExecStart=/bin/rk_conexoesWifi

[Install]
WantedBy=multi-user.target
